# src/crontab
0       2       *       *       *       python /app/deletion.py