# src/crontab
*/35    *       *       *       *       python /app/convert.py
0       2       *       *       *       python /app/deletion.py